# blog/blog.service
[Unit]
Description=Blog Node.js Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/blog/api
EnvironmentFile=/var/www/blog/api/.env
ExecStart=/usr/bin/node /var/www/blog/api/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blog

[Install]
WantedBy=multi-user.target
